# etc/cron.d/tms
# Update System
*/15 * * * *	root            /usr/local/sbin/tms_update >> /dev/null 2>&1

